# backup/inotifywait-backup.service
[Service]
DynamicUser=true
ExecStart=bash /path/to/inotify.sh

[Install]
WantedBy=multi-user.target
